# cron.weekly/certbot
#!/bin/sh -e
/opt/certbot/certbot-auto renew
nginx -s reload

